# supervisor/gruffalo.conf
[program:gruffalo]
command=/gruffalo/gruffalo-standalone.sh
stdout_logfile=/var/log/gruffalo.log
directory=/gruffalo
user=root
autostart=true
autorestart=true
redirect_stderr=True
stopasgroup=true
